# docker/entrails/nginx/supervisor/supervisor-tunnel-log.conf
[program:twilio-tunnel-log]
command=/ngrok/ngrok_tunnel_logger.sh
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
